# cron/lis.run.cron
#!/bin/csh
set finddate = finddate.sh
set homedir=/gpfs/dell2/emc/retros/noscrub/Youlong.Xia/
set yyyy = `date +%Y`
set   mm = `date +%m`
set   dd = `date +%d`
set yyyymmdd2 = $yyyy$mm$dd
set yyyymmdd1 = `sh $finddate $yyyymmdd2 d-3`
$homedir/gldas.v2.3.0/scripts/lis.run.sh $yyyymmdd1 $yyyymmdd2

